# deploy/bugout-box.service
[Unit]
Description=bug-out box voice assistant
After=network.target ollama.service
Wants=ollama.service

[Service]
Type=simple
User=caleb
WorkingDirectory=/home/caleb/bugout-box
Environment=BOX_VAULT=/media/caleb/Expansion/vault
ExecStart=/usr/bin/python3 /home/caleb/bugout-box/run.py
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
