# deploy/ember.service
[Unit]
Description=Ember brain — wake-word voice loop (STT/RAG/Gemma/TTS)
After=ollama.service ember-whisper.service sound.target

[Service]
User=caleb
WorkingDirectory=/home/caleb/bugout-box
Environment=PYTHONUNBUFFERED=1
ExecStart=/usr/bin/python3 /home/caleb/bugout-box/run.py
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
